# launchd/com.claude-vault.watchdog.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!-- memvault watchdog — install.sh renders __VAULT_BIN__ and drops this in ~/Library/LaunchAgents/ -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.memvault.watchdog</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__VAULT_BIN__/watchdog.sh</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>CLAUDE_VAULT_CONFIG</key>
        <string>__VAULT_CONFIG__</string>
    </dict>
    <!-- every 3h + at login -->
    <key>StartInterval</key>
    <integer>10800</integer>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
